# loadwatch.cron
#!/bin/bash
*/3 * * * * root /usr/bin/flock -n /dev/shm/loadwatch.lock -c /usr/local/lp/bin/loadwatch &> /dev/null
0   1 * * * root source /etc/loadwatch.conf; /usr/bin/find ${DIR} -maxdepth 1 -type f -mtime +30 -delete
0   2 * * 0 root source /etc/loadwatch.conf; /usr/bin/tail -n100000 ${DIR}.log > ${DIR}.new; mv -f ${DIR}.new ${DIR}.log